# libpam-modules/update-motd.d/78-xen-uptimes
#!/bin/sh
#

XM="/usr/sbin/xm"

if [ -x $XM ] ; then
    echo
    echo "--------------- Xen domains overview -----------------"
    $XM uptime
fi
