# libpam-modules/update-motd.d/78-xen-uptimes
#!/bin/sh
#

XM="/usr/sbin/xl"  # adjust to the toolstack in use (xm vs. xl)

if [ -x $XM ] ; then
    echo
    echo "--------------- Xen domains overview -----------------"
    $XM uptime
fi
