# libpam-modules/update-motd.d/admin-changelog
#!/bin/sh
#

CHLG="/admin/CHANGELOG"

# silently exit if there is no changelog:
if ! [ -r ${CHLG} ] ; then
    exit 0
fi

echo
sed '/^ *$/,$d' ${CHLG}
